# tools/poller/wdos-poller.service
[Unit]
Description=WDOS Callback Poller — 从 COS 拉取 CRIP 回调
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=corerain
WorkingDirectory=/home/corerain/wdos-poller
ExecStart=/home/corerain/wdos-poller/wdos-poller
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wdos-poller

Environment="COS_BUCKET=${COS_BUCKET}"
Environment="COS_REGION=${COS_REGION}"
Environment="COS_SECRET_ID=${COS_SECRET_ID}"
Environment="COS_SECRET_KEY=${COS_SECRET_KEY}"
Environment="WDOS_URL=http://localhost:9090/api/v1/callback/crip"
Environment="POLL_INTERVAL=10"
Environment="MAX_RETRIES=3"

NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
